# application/supervisord.conf
[program:node]
command = node app-dotcloud.js
directory = /home/dotcloud/current
